# app/cron/crontab
* * * * * cd /var/www/lotti && php artisan schedule:run >> /var/log/cron.log 2>&1
